# deploy/supervisor.conf
[program:order_consumer]
command=sudo -u www-data php /var/www/demo/current/artisan rabbitmq:consume --queue=order
process_name=order_consumer_%(process_num)02d
numprocs=1
directory=/tmp
autostart=true
autorestart=true
startsecs=3
startretries=10
user=www-data
redirect_stderr=false
stdout_logfile=/var/www/log/supervisor.order_consumer.out.log
stdout_capture_maxbytes=1MB
stderr_logfile=/var/www/log/supervisor.order_consumer.error.log
stderr_capture_maxbytes=1MB