# hw17/docker/workers/supervisor.conf
[supervisord]
logfile = /tmp/supervisord.log
logfile_maxbytes = 50MB
logfile_backups=10
loglevel = info
pidfile = /tmp/supervisord.pid
nodaemon = true
minfds = 1024
minprocs = 200
umask = 022
user = root
identifier = supervisor
directory = /tmp
nocleanup = true
childlogdir = /tmp
strip_ansi = false

[program:php_worker]
command=php /var/www/worker.php
numprocs=2
process_name=%(program_name)s_%(process_num)02d
directory=/tmp
autostart=true
autorestart=true
startsecs=5
startretries=10
redirect_stderr=false
stdout_logfile=/var/log/php_workers.out.log
stdout_capture_maxbytes=1MB
stderr_logfile=/var/log/php_workers.error.log
stderr_capture_maxbytes=1MB